# notebooks/crontab
* * * * * python /notebooks/model_server/check_s3.py >> /var/log/cron.log 2>&1
